# data/systemd/bodhi-update-manager-refresh.timer
[Unit]
Description=Update Manager — hourly APT cache refresh timer
Documentation=https://github.com/BodhiLinux/bodhi-update-manager

[Timer]
# First activation: 10 minutes after boot/timer-start.
OnBootSec=10min
# Subsequent activations: every hour.
OnUnitActiveSec=1h
# Spread activations randomly within a 5-minute window to avoid herd effect.
RandomizedDelaySec=5min
# Run a missed activation on the next boot (e.g. laptop was suspended).
Persistent=true
Unit=bodhi-update-manager-refresh.service

[Install]
WantedBy=timers.target
